# templates/mosquitto/systemd/mosquitto-podman.service.j2
{{ ansible_managed | comment }}
# container-{{ mosquitto_id }}.service
# autogenerated by Podman 1.6.4

[Unit]
Description=Podman container-{{ mosquitto_id }}
Documentation=man:podman-generate-systemd(1)

[Service]
Restart=on-failure
ExecStart=/usr/bin/podman start {{ mosquitto_id }}
ExecStop=/usr/bin/podman stop -t 10 {{ mosquitto_id }}
KillMode=none
Type=forking
PIDFile=/var/run/containers/storage/overlay-containers/{{ mosquitto_id }}/userdata/conmon.pid

[Install]
WantedBy=multi-user.target
